# packaging/linux/specter-dns.service
[Unit]
Description=Specter ACME DNS solver
After=network-online.target specter-server.service
Wants=network-online.target
ConditionPathExists=/etc/specter/dns.env

[Service]
Type=simple
User=specter
Group=specter
RuntimeDirectory=specter
RuntimeDirectoryMode=0770
ExecStart=/usr/libexec/specter/dns-launch
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=5s
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
